# kage.service.in
[Unit]
Description=Kage Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=INSTALL_USER
ExecStart=SBINDIR/kageserver SYSCONFDIR/kage.cfg
StandardOutput=append:LOCALSTATEDIR/log/kage.log

[Install]
WantedBy=multi-user.target
