# kage.service.in
[Unit]
Description=Kage Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=INSTALL_USER
ExecStart=SBINDIR/kageserver SYSCONFDIR/kage.cfg
StandardOutput=append:/var/log/kage/kage.log
# If dumping net data
#WorkingDirectory=/var/lib/kage

[Install]
WantedBy=multi-user.target
